# supervisord.conf
[supervisord]
nodaemon=true

[program:npm]
command=serve -s info_quest_frontend/build

[program:python]
command=bash info_quest_backend/run.sh
