# supervisord.conf
[supervisord]
nodaemon=true

[program:npm]
command=serve -s frontend/build

[program:python]
command=bash ./backend/run.sh
